# check_rippled.conf
#!upstart
description "Check Rippled"

start on (local-filesystems and net-device-up IFACE!=lo)
stop on [!12345]

console log

setuid afrank
setgid afrank
exec /home/afrank/check_rippled_daemon.py

